# ports-mgmt/packagekit/pkg-plist
bin/pkcon
bin/pkmon
@sample %%ETCDIR%%/CommandNotFound.conf.default %%ETCDIR%%/CommandNotFound.conf
@sample %%ETCDIR%%/PackageKit.conf.default %%ETCDIR%%/PackageKit.conf
%%ETCDIR%%/Vendor.conf
etc/cron.daily/packagekit-background.cron
etc/dbus-1/system.d/org.freedesktop.PackageKit.conf
etc/profile.d/PackageKit.sh
etc/sysconfig/packagekit-background
include/PackageKit/packagekit-glib2/packagekit.h
include/PackageKit/packagekit-glib2/pk-bitfield.h
include/PackageKit/packagekit-glib2/pk-category.h
include/PackageKit/packagekit-glib2/pk-client-helper.h
include/PackageKit/packagekit-glib2/pk-client-sync.h
include/PackageKit/packagekit-glib2/pk-client.h
include/PackageKit/packagekit-glib2/pk-common.h
include/PackageKit/packagekit-glib2/pk-control-sync.h
include/PackageKit/packagekit-glib2/pk-control.h
include/PackageKit/packagekit-glib2/pk-desktop.h
include/PackageKit/packagekit-glib2/pk-details.h
include/PackageKit/packagekit-glib2/pk-distro-upgrade.h
include/PackageKit/packagekit-glib2/pk-enum-types.h
include/PackageKit/packagekit-glib2/pk-enum.h
include/PackageKit/packagekit-glib2/pk-error.h
include/PackageKit/packagekit-glib2/pk-eula-required.h
include/PackageKit/packagekit-glib2/pk-files.h
include/PackageKit/packagekit-glib2/pk-item-progress.h
include/PackageKit/packagekit-glib2/pk-media-change-required.h
include/PackageKit/packagekit-glib2/pk-offline.h
include/PackageKit/packagekit-glib2/pk-package-id.h
include/PackageKit/packagekit-glib2/pk-package-ids.h
include/PackageKit/packagekit-glib2/pk-package-sack-sync.h
include/PackageKit/packagekit-glib2/pk-package-sack.h
include/PackageKit/packagekit-glib2/pk-package.h
include/PackageKit/packagekit-glib2/pk-progress.h
include/PackageKit/packagekit-glib2/pk-repo-detail.h
include/PackageKit/packagekit-glib2/pk-repo-signature-required.h
include/PackageKit/packagekit-glib2/pk-require-restart.h
include/PackageKit/packagekit-glib2/pk-results.h
include/PackageKit/packagekit-glib2/pk-source.h
include/PackageKit/packagekit-glib2/pk-task-sync.h
include/PackageKit/packagekit-glib2/pk-task.h
include/PackageKit/packagekit-glib2/pk-transaction-list.h
include/PackageKit/packagekit-glib2/pk-transaction-past.h
include/PackageKit/packagekit-glib2/pk-update-detail.h
include/PackageKit/packagekit-glib2/pk-version.h
lib/girepository-1.0/PackageKitGlib-1.0.typelib
lib/libpackagekit-glib2.so
lib/libpackagekit-glib2.so.18
lib/libpackagekit-glib2.so.18.1.3
lib/packagekit-backend/libpk_backend_dummy.so
%%PORTS%%lib/packagekit-backend/libpk_backend_ports.so
lib/packagekit-backend/libpk_backend_test_fail.so
lib/packagekit-backend/libpk_backend_test_nop.so
lib/packagekit-backend/libpk_backend_test_spawn.so
lib/packagekit-backend/libpk_backend_test_succeed.so
lib/packagekit-backend/libpk_backend_test_thread.so
libdata/pkgconfig/packagekit-glib2.pc
libexec/packagekit-direct
libexec/packagekitd
libexec/pk-command-not-found
%%MAN%%man/man1/pkcon.1.gz
%%MAN%%man/man1/pkmon.1.gz
%%PORTS%%%%DATADIR%%/helpers/ports/portsBackend.rb
%%PORTS%%%%DATADIR%%/helpers/ports/ruby_packagekit/backend.rb
%%PORTS%%%%DATADIR%%/helpers/ports/ruby_packagekit/enums.rb
%%PORTS%%%%DATADIR%%/helpers/ports/ruby_packagekit/prints.rb
%%DATADIR%%/helpers/test_spawn/search-name.sh
%%DATADIR%%/pk-upgrade-distro.sh
%%DATADIR%%/transactions.db
@postexec if [ ! -f /var/lib/PackageKit/%f ]; then cp -f %B/%f /var/lib/PackageKit; chmod o-r /var/lib/PackageKit/%f; fi
%%BASH%%share/bash-completion/completions/pkcon
share/dbus-1/interfaces/org.freedesktop.PackageKit.Transaction.xml
share/dbus-1/interfaces/org.freedesktop.PackageKit.xml
share/dbus-1/system-services/org.freedesktop.PackageKit.service
share/gir-1.0/PackageKitGlib-1.0.gir
%%DOCS%%share/gtk-doc/html/PackageKit/PackageKit-Common-functions.html
%%DOCS%%share/gtk-doc/html/PackageKit/PackageKit-Enumerations.html
%%DOCS%%share/gtk-doc/html/PackageKit/PackageKit-Offline-Updates.html
%%DOCS%%share/gtk-doc/html/PackageKit/PackageKit-PackageIDs.html
%%DOCS%%share/gtk-doc/html/PackageKit/PackageKit-PkBitfield.html
%%DOCS%%share/gtk-doc/html/PackageKit/PackageKit-PkCategory.html
%%DOCS%%share/gtk-doc/html/PackageKit/PackageKit-PkClient.html
%%DOCS%%share/gtk-doc/html/PackageKit/PackageKit-PkClientHelper.html
%%DOCS%%share/gtk-doc/html/PackageKit/PackageKit-PkControl.html
%%DOCS%%share/gtk-doc/html/PackageKit/PackageKit-PkDesktop.html
%%DOCS%%share/gtk-doc/html/PackageKit/PackageKit-PkDetails.html
%%DOCS%%share/gtk-doc/html/PackageKit/PackageKit-PkDistroUpgrade.html
%%DOCS%%share/gtk-doc/html/PackageKit/PackageKit-PkError.html
%%DOCS%%share/gtk-doc/html/PackageKit/PackageKit-PkEulaRequired.html
%%DOCS%%share/gtk-doc/html/PackageKit/PackageKit-PkFiles.html
%%DOCS%%share/gtk-doc/html/PackageKit/PackageKit-PkItemProgress.html
%%DOCS%%share/gtk-doc/html/PackageKit/PackageKit-PkMediaChangeRequired.html
%%DOCS%%share/gtk-doc/html/PackageKit/PackageKit-PkPackage.html
%%DOCS%%share/gtk-doc/html/PackageKit/PackageKit-PkPackageSack.html
%%DOCS%%share/gtk-doc/html/PackageKit/PackageKit-PkProgress.html
%%DOCS%%share/gtk-doc/html/PackageKit/PackageKit-PkRepoDetail.html
%%DOCS%%share/gtk-doc/html/PackageKit/PackageKit-PkRepoSignatureRequired.html
%%DOCS%%share/gtk-doc/html/PackageKit/PackageKit-PkRequireRestart.html
%%DOCS%%share/gtk-doc/html/PackageKit/PackageKit-PkResults.html
%%DOCS%%share/gtk-doc/html/PackageKit/PackageKit-PkSource.html
%%DOCS%%share/gtk-doc/html/PackageKit/PackageKit-PkTask.html
%%DOCS%%share/gtk-doc/html/PackageKit/PackageKit-PkTransactionList.html
%%DOCS%%share/gtk-doc/html/PackageKit/PackageKit-PkTransactionPast.html
%%DOCS%%share/gtk-doc/html/PackageKit/PackageKit-PkUpdateDetail.html
%%DOCS%%share/gtk-doc/html/PackageKit/PackageKit.devhelp2
%%DOCS%%share/gtk-doc/html/PackageKit/PackageKit.html
%%DOCS%%share/gtk-doc/html/PackageKit/Transaction.html
%%DOCS%%share/gtk-doc/html/PackageKit/api-index-0-5-0.html
%%DOCS%%share/gtk-doc/html/PackageKit/api-index-0-5-2.html
%%DOCS%%share/gtk-doc/html/PackageKit/api-index-0-5-3.html
%%DOCS%%share/gtk-doc/html/PackageKit/api-index-0-5-4.html
%%DOCS%%share/gtk-doc/html/PackageKit/api-index-0-5-5.html
%%DOCS%%share/gtk-doc/html/PackageKit/api-index-0-6-0.html
%%DOCS%%share/gtk-doc/html/PackageKit/api-index-0-6-1.html
%%DOCS%%share/gtk-doc/html/PackageKit/api-index-0-6-10.html
%%DOCS%%share/gtk-doc/html/PackageKit/api-index-0-6-11.html
%%DOCS%%share/gtk-doc/html/PackageKit/api-index-0-6-13.html
%%DOCS%%share/gtk-doc/html/PackageKit/api-index-0-6-2.html
%%DOCS%%share/gtk-doc/html/PackageKit/api-index-0-6-3.html
%%DOCS%%share/gtk-doc/html/PackageKit/api-index-0-6-4.html
%%DOCS%%share/gtk-doc/html/PackageKit/api-index-0-6-5.html
%%DOCS%%share/gtk-doc/html/PackageKit/api-index-0-7-2.html
%%DOCS%%share/gtk-doc/html/PackageKit/api-index-0-7-5.html
%%DOCS%%share/gtk-doc/html/PackageKit/api-index-0-8-0.html
%%DOCS%%share/gtk-doc/html/PackageKit/api-index-0-8-1.html
%%DOCS%%share/gtk-doc/html/PackageKit/api-index-0-8-11.html
%%DOCS%%share/gtk-doc/html/PackageKit/api-index-0-8-12.html
%%DOCS%%share/gtk-doc/html/PackageKit/api-index-0-8-14.html
%%DOCS%%share/gtk-doc/html/PackageKit/api-index-0-8-16.html
%%DOCS%%share/gtk-doc/html/PackageKit/api-index-0-8-17.html
%%DOCS%%share/gtk-doc/html/PackageKit/api-index-0-8-2.html
%%DOCS%%share/gtk-doc/html/PackageKit/api-index-0-8-6.html
%%DOCS%%share/gtk-doc/html/PackageKit/api-index-0-8-8.html
%%DOCS%%share/gtk-doc/html/PackageKit/api-index-0-9-1.html
%%DOCS%%share/gtk-doc/html/PackageKit/api-index-0-9-5.html
%%DOCS%%share/gtk-doc/html/PackageKit/api-index-0-9-6.html
%%DOCS%%share/gtk-doc/html/PackageKit/api-index-1-0-10.html
%%DOCS%%share/gtk-doc/html/PackageKit/api-index-1-0-12.html
%%DOCS%%share/gtk-doc/html/PackageKit/api-index-1-0-2.html
%%DOCS%%share/gtk-doc/html/PackageKit/api-index-1-1-2.html
%%DOCS%%share/gtk-doc/html/PackageKit/api-index-1-1-8.html
%%DOCS%%share/gtk-doc/html/PackageKit/api-reference.html
%%DOCS%%share/gtk-doc/html/PackageKit/backend-compiled.html
%%DOCS%%share/gtk-doc/html/PackageKit/backend-spawn.html
%%DOCS%%share/gtk-doc/html/PackageKit/concepts.html
%%DOCS%%share/gtk-doc/html/PackageKit/home.png
%%DOCS%%share/gtk-doc/html/PackageKit/index.html
%%DOCS%%share/gtk-doc/html/PackageKit/introduction-backends.html
%%DOCS%%share/gtk-doc/html/PackageKit/introduction-cancellation.html
%%DOCS%%share/gtk-doc/html/PackageKit/introduction-errors.html
%%DOCS%%share/gtk-doc/html/PackageKit/introduction-group-type.html
%%DOCS%%share/gtk-doc/html/PackageKit/introduction-ideas-filters.html
%%DOCS%%share/gtk-doc/html/PackageKit/introduction-ideas-status.html
%%DOCS%%share/gtk-doc/html/PackageKit/introduction-ideas-transactionid.html
%%DOCS%%share/gtk-doc/html/PackageKit/introduction-ideas-transactions.html
%%DOCS%%share/gtk-doc/html/PackageKit/introduction.html
%%DOCS%%share/gtk-doc/html/PackageKit/left-insensitive.png
%%DOCS%%share/gtk-doc/html/PackageKit/left.png
%%DOCS%%share/gtk-doc/html/PackageKit/lpackagekit-glib2.html
%%DOCS%%share/gtk-doc/html/PackageKit/pk-structure.png
%%DOCS%%share/gtk-doc/html/PackageKit/pk-structure.svg
%%DOCS%%share/gtk-doc/html/PackageKit/pk-transactions-auto-untrusted.png
%%DOCS%%share/gtk-doc/html/PackageKit/pk-transactions-download.png
%%DOCS%%share/gtk-doc/html/PackageKit/pk-transactions-failure.png
%%DOCS%%share/gtk-doc/html/PackageKit/pk-transactions-repair-required.png
%%DOCS%%share/gtk-doc/html/PackageKit/pk-transactions-set-locale.png
%%DOCS%%share/gtk-doc/html/PackageKit/pk-transactions-sig-install.png
%%DOCS%%share/gtk-doc/html/PackageKit/pk-transactions-success.png
%%DOCS%%share/gtk-doc/html/PackageKit/pk-transactions-trusted.png
%%DOCS%%share/gtk-doc/html/PackageKit/pk-transactions.svg
%%DOCS%%share/gtk-doc/html/PackageKit/right-insensitive.png
%%DOCS%%share/gtk-doc/html/PackageKit/right.png
%%DOCS%%share/gtk-doc/html/PackageKit/specification.html
%%DOCS%%share/gtk-doc/html/PackageKit/style.css
%%DOCS%%share/gtk-doc/html/PackageKit/up-insensitive.png
%%DOCS%%share/gtk-doc/html/PackageKit/up.png
%%NLS%%share/locale/ar/LC_MESSAGES/PackageKit.mo
%%NLS%%share/locale/as/LC_MESSAGES/PackageKit.mo
%%NLS%%share/locale/ast/LC_MESSAGES/PackageKit.mo
%%NLS%%share/locale/az/LC_MESSAGES/PackageKit.mo
%%NLS%%share/locale/bg/LC_MESSAGES/PackageKit.mo
%%NLS%%share/locale/bn_IN/LC_MESSAGES/PackageKit.mo
%%NLS%%share/locale/br/LC_MESSAGES/PackageKit.mo
%%NLS%%share/locale/ca/LC_MESSAGES/PackageKit.mo
%%NLS%%share/locale/ca@valencia/LC_MESSAGES/PackageKit.mo
%%NLS%%share/locale/cs/LC_MESSAGES/PackageKit.mo
%%NLS%%share/locale/cy/LC_MESSAGES/PackageKit.mo
%%NLS%%share/locale/da/LC_MESSAGES/PackageKit.mo
%%NLS%%share/locale/de/LC_MESSAGES/PackageKit.mo
%%NLS%%share/locale/el/LC_MESSAGES/PackageKit.mo
%%NLS%%share/locale/en_GB/LC_MESSAGES/PackageKit.mo
%%NLS%%share/locale/eo/LC_MESSAGES/PackageKit.mo
%%NLS%%share/locale/es/LC_MESSAGES/PackageKit.mo
%%NLS%%share/locale/eu/LC_MESSAGES/PackageKit.mo
%%NLS%%share/locale/fa/LC_MESSAGES/PackageKit.mo
%%NLS%%share/locale/fi/LC_MESSAGES/PackageKit.mo
%%NLS%%share/locale/fo/LC_MESSAGES/PackageKit.mo
%%NLS%%share/locale/fr/LC_MESSAGES/PackageKit.mo
%%NLS%%share/locale/ga/LC_MESSAGES/PackageKit.mo
%%NLS%%share/locale/gl/LC_MESSAGES/PackageKit.mo
%%NLS%%share/locale/gu/LC_MESSAGES/PackageKit.mo
%%NLS%%share/locale/he/LC_MESSAGES/PackageKit.mo
%%NLS%%share/locale/hi/LC_MESSAGES/PackageKit.mo
%%NLS%%share/locale/hr/LC_MESSAGES/PackageKit.mo
%%NLS%%share/locale/hu/LC_MESSAGES/PackageKit.mo
%%NLS%%share/locale/ia/LC_MESSAGES/PackageKit.mo
%%NLS%%share/locale/id/LC_MESSAGES/PackageKit.mo
%%NLS%%share/locale/it/LC_MESSAGES/PackageKit.mo
%%NLS%%share/locale/ja/LC_MESSAGES/PackageKit.mo
%%NLS%%share/locale/ka/LC_MESSAGES/PackageKit.mo
%%NLS%%share/locale/kk/LC_MESSAGES/PackageKit.mo
%%NLS%%share/locale/kn/LC_MESSAGES/PackageKit.mo
%%NLS%%share/locale/ko/LC_MESSAGES/PackageKit.mo
%%NLS%%share/locale/lt/LC_MESSAGES/PackageKit.mo
%%NLS%%share/locale/lv/LC_MESSAGES/PackageKit.mo
%%NLS%%share/locale/ml/LC_MESSAGES/PackageKit.mo
%%NLS%%share/locale/mr/LC_MESSAGES/PackageKit.mo
%%NLS%%share/locale/ms/LC_MESSAGES/PackageKit.mo
%%NLS%%share/locale/nb/LC_MESSAGES/PackageKit.mo
%%NLS%%share/locale/nl/LC_MESSAGES/PackageKit.mo
%%NLS%%share/locale/nn/LC_MESSAGES/PackageKit.mo
%%NLS%%share/locale/oc/LC_MESSAGES/PackageKit.mo
%%NLS%%share/locale/or/LC_MESSAGES/PackageKit.mo
%%NLS%%share/locale/pa/LC_MESSAGES/PackageKit.mo
%%NLS%%share/locale/pl/LC_MESSAGES/PackageKit.mo
%%NLS%%share/locale/pt/LC_MESSAGES/PackageKit.mo
%%NLS%%share/locale/pt_BR/LC_MESSAGES/PackageKit.mo
%%NLS%%share/locale/ro/LC_MESSAGES/PackageKit.mo
%%NLS%%share/locale/ru/LC_MESSAGES/PackageKit.mo
%%NLS%%share/locale/sk/LC_MESSAGES/PackageKit.mo
%%NLS%%share/locale/sl/LC_MESSAGES/PackageKit.mo
%%NLS%%share/locale/sq/LC_MESSAGES/PackageKit.mo
%%NLS%%share/locale/sr/LC_MESSAGES/PackageKit.mo
%%NLS%%share/locale/sr@latin/LC_MESSAGES/PackageKit.mo
%%NLS%%share/locale/sv/LC_MESSAGES/PackageKit.mo
%%NLS%%share/locale/ta/LC_MESSAGES/PackageKit.mo
%%NLS%%share/locale/te/LC_MESSAGES/PackageKit.mo
%%NLS%%share/locale/th/LC_MESSAGES/PackageKit.mo
%%NLS%%share/locale/tr/LC_MESSAGES/PackageKit.mo
%%NLS%%share/locale/uk/LC_MESSAGES/PackageKit.mo
%%NLS%%share/locale/vi/LC_MESSAGES/PackageKit.mo
%%NLS%%share/locale/wa/LC_MESSAGES/PackageKit.mo
%%NLS%%share/locale/zh_CN/LC_MESSAGES/PackageKit.mo
%%NLS%%share/locale/zh_HK/LC_MESSAGES/PackageKit.mo
%%NLS%%share/locale/zh_TW/LC_MESSAGES/PackageKit.mo
share/polkit-1/actions/org.freedesktop.packagekit.policy
share/polkit-1/rules.d/org.freedesktop.packagekit.rules
%%VAPI%%share/vala/vapi/packagekit-glib2.vapi
@dir /var/cache/PackageKit/downloads
@dir /var/cache/PackageKit/metadata
@dir /var/cache/PackageKit
@dir /var/lib/PackageKit
@dir /var/lib
